# pi-config/systemd/ptt.service
# Installed to /etc/systemd/system/ptt.service during explicit office/dev
# target setup. Ordinary deploy.sh releases do not alter machine config. This is the
# legacy recplay/ptt Pi runtime; deploy.sh does not install kibod.
# This process owns its local flat turns.jsonl; it does not upload to or
# coordinate authority with kibod.
# Logs: journalctl -u ptt -f
[Unit]
Description=kibo push-to-talk (joystick record/playback loop)
After=sound.target

[Service]
User=jesse
Environment=HOME=/home/jesse
# GEMINI_API_KEY lives here; target setup installs it explicitly and releases
# preserve it.
EnvironmentFile=-/home/jesse/.env
WorkingDirectory=/home/jesse
ExecStart=/home/jesse/ptt
# Also covers the joystick being unplugged at boot: keep retrying until
# it appears.
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
